# PI3Beacon/beacon.service
[Unit]
Description=BLE Beacon

[Service]
Type=simple
ExecStart=/usr/bin/beacon_final.py

[Install]
WantedBy=multi-user.target
